# build_files/files/systemd/user/post-install.service
[Unit]
Description=Run user post-install tasks after login (dotfiles, brew, Doom, etc.)
After=graphical-session.target network-online.target
Wants=graphical-session.target network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/post-install-user.sh
# ensure we have a usable PATH for user tools like brew/doom
Environment=PATH=%h/.local/bin:%h/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/bin:/usr/bin:/usr/local/bin

[Install]
WantedBy=default.target
